# deploy/dns/domain-ddns-update.service
[Unit]
Description="Update ddns of playreadingparty.com"

[Service]
ExecStart=/bin/domain-ddns-update.sh
