# deploy/dns/domain-ddns-update.timer
[Unit]
Description=Run ddns update for playreadingparty.com

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
